# build/scripts/vm/valekos-vm-detect.service
[Unit]
Description=ValekOS VM Detection and Optimization
After=multi-user.target
Before=display-manager.service
ConditionPathExists=!/run/valekos/vm-info

[Service]
Type=oneshot
ExecStart=/usr/bin/valekos-vm-detect
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=multi-user.target
